# layer/meta-rapid/recipes-apps/rapid/rapid-headless-init/rapid-headless.service
[Unit]
Description=Rapid Headless Timer Service
After=network.target

[Service]
ExecStart=/usr/bin/rapid_headless --gpsd
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
